# gategpt/rootfs/etc/services.d/gategpt/finish
#!/usr/bin/with-contenv bash
# Called by s6 when the service receives a SIGTERM/SIGINT
# or when the container is stopping.

set -e

echo "[gategpt] Shutting down – cleaning up…"

# 1. Tell the Node process to exit cleanly so that
#    whatsapp-web.js can close its WebSocket and save state.
if [ -n "${S6_SERVICE_PID-}" ]; then
  kill -s SIGINT "$S6_SERVICE_PID" 2>/dev/null || true
  # Give it a moment (max 10 s) before forcing it.
  for _ in $(seq 1 10); do
    kill -0 "$S6_SERVICE_PID" 2>/dev/null || break
    sleep 1
  done
fi

# 2. House-keeping (optional).  For example, remove the last QR code
rm -f /share/gategpt/qr.png || true

echo "[gategpt] Cleanup complete."
exit 0