# gategpt/rootfs/etc/services.d/gategpt/run
#!/usr/bin/env bash
set -e

# Export every add-on option as ENV for Node
while IFS='=' read -r name value ; do
  export "$name"="$value"
done < <(env | grep -E '^(GATE_|OPENAI_|PUSHOVER_|CHATGPT_|MESSAGE_|IGNORE_FILE|TRIGGER_|.*_DELAY_MS|MAX_MESSAGES_)')

# Ensure a place to drop the QR code
mkdir -p /share/gategpt
export QR_PATH=/share/gategpt/qr.png

exec node main.js